# node-lukeburroughs.conf
#!upstart

env APP_FILE=/home/node-lukeburroughs/angularjs/app.js
env APP_NAME=lukeburroughs
env APP_USER=node-lukeburroughs
env APP_PATH=/home/$APP_USER

description "Luke Burroughs portfolio node service"

start on filesystem or runlevel [2345]
stop on shutdown

expect fork

respawn
respawn limit 5 30

script
  exec su - $APP_USER -c forever start -a -l /tmp/forever.log -o /tmp/forever.$APP_NAME.stdout.log -e /tmp/forever.$APP_NAME.stderr.log $APP_FILE
end script

pre-stop script
  exec forever stop $APP_SERVER_FILE
end script